# restic-backup.service
[Unit]
Description=Restic Backup Service

[Service]
Type=oneshot
ExecStart=/home/alialre/scripts/restic-backup.sh
